# server/config/crontab
*/30 * * * * python /secret-club/src/main.py
