# server/config/crontab
*/30 * * * * python3 ~/tmd/main.py >> ~/tmd/log.txt
